# deploy/docker/develop/consumer/supervisor.conf
[program:outbox-process]
command=/app/cmd outbox:process
numprocs=1
autostart=true
startsecs=10
autorestart=true
startretries=17280
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:outbox-prune]
command=/app/cmd outbox:prune --remove-older-than-hour 24 --run-every-hour 1
numprocs=1
autostart=true
startsecs=10
autorestart=true
startretries=17280
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0


[program:create_ozone_card_segment_consumer]
command=/app/cmd consumer --create-ozone-card-segment
numprocs=1
autostart=true
startsecs=10
autorestart=true
startretries=17280
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:user-consumer]
command=/app/cmd consumer --user-consumer
numprocs=1
autostart=true
startsecs=10
autorestart=true
startretries=17280
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:okcs-product-sync-consumer]
command=/app/cmd consumer --okcs-product-sync
numprocs=1
autostart=true
startsecs=10
autorestart=true
startretries=17280
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:okcs-category-sync-consumer]
command=/app/cmd consumer --okcs-category-sync
numprocs=1
autostart=true
startsecs=10
autorestart=true
startretries=17280
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0


[program:okcs-brand-sync-consumer]
command=/app/cmd consumer --okcs-brand-sync
numprocs=1
autostart=true
startsecs=10
autorestart=true
startretries=17280
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:okcs-city-sync-consumer]
command=/app/cmd consumer --okcs-city-sync
numprocs=1
autostart=true
startsecs=10
autorestart=true
startretries=17280
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:okcs-branch-sync-consumer]
command=/app/cmd consumer --okcs-branch-sync
numprocs=1
autostart=true
startsecs=10
autorestart=true
startretries=17280
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:okcs-branch-product-initial-sync]
command=/app/cmd consumer --okcs-branch-product-initial-sync
numprocs=1
autostart=true
startsecs=10
autorestart=true
startretries=17280
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:okcs-branch-product-sync]
command=/app/cmd consumer --okcs-branch-product-sync
numprocs=1
autostart=true
startsecs=10
autorestart=true
startretries=17280
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:okcs-branch-product-batch-upsert]
command=/app/cmd consumer --okcs-branch-product-batch-upsert
numprocs=1
autostart=true
startsecs=10
autorestart=true
startretries=17280
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:okcs-sync-scheduler]
command=/app/cmd okcs-sync --scheduler
numprocs=1
autostart=true
startsecs=10
autorestart=true
startretries=17280
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0


[program:assign-voucher-by-item-consumer]
command=/app/cmd consumer --assign-voucher-by-item
numprocs=1
autostart=true
startsecs=10
autorestart=true
startretries=17280
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:create-voucher-by-item-consumer]
command=/app/cmd consumer --create-voucher-by-item
numprocs=1
autostart=true
startsecs=10
autorestart=true
startretries=17280
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0